# init/alpinelinux/socat-forward-service
#!/sbin/openrc-run

name="socat-forward"
description="One-shot socat forward starter"

command="/usr/local/socat-forward/socat-starter.sh"
